# qwen_bot.service
[Unit]
Description=Qwen Exam Assistant Telegram Bot
After=network.target

[Service]
User=lh
Group=lh
WorkingDirectory=/opt/ococl/qwen_bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/ococl/qwen_bot/.venv/bin/python3 main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qwen_bot

[Install]
WantedBy=multi-user.target
